# Scripts/mountpoint/media-naslocal-music.automount
[Unit]
Description=NAS SMB db1 automount

[Automount]
Where=/media/naslocal/music
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
