# Scripts/mountpoint/media-nasremote-music.automount
[Unit]
Description=NAS SMB music automount

[Automount]
Where=/media/nasremote/music
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
