# Scripts/mountpoint/media-nasremote-db1.automount
[Unit]
Description=NAS SMB db1 automount

[Automount]
Where=/media/nasremote/db1
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
